# debian/cyborgbackup.cyborgbackup-celery.service
[Unit]
Description=CyBorgBackup Celery Service
PartOf=cyborgbackup.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/cyborgbackup
WorkingDirectory=/opt/cyborgbackup
User=cyborgbackup
Group=cyborgbackup
ExecStart=/opt/cyborgbackup/venv/bin/celery -A cyborgbackup worker -E -c 1

[Install]
WantedBy=multi-user.target